# rootfs.list
#dir		<name>																			<mode>	<uid>	<gid>
dir		/sbin																				0755		0		0
dir		/bin																				0755		0		0
	
#file		<name>					<location>											<mode>	<uid>	<gid>	[<hard links>]
file		/bin/busybox			/tks/dctab/br_out-skel/bin/busybox			4755		0		0
file		/etc/inittab			/tks/dctab/etc/inittab							0644		0		0
file		/etc/fstab				/tks/dctab/etc/fstab								0644		0		0
file		/etc/init.d/S01dctab	/tks/dctab/etc/init.d/S01dctab				0755		0		0

#slink	<name>					<target>												<mode>	<uid>	<gid>
slink		/sbin/fsck				../bin/busybox										777		0		0
slink		/sbin/getty		   	../bin/busybox										777		0		0
slink		/sbin/init				../bin/busybox										777		0		0
slink		/bin/ash					busybox												777		0		0
slink		/bin/ln					busybox												777		0		0
slink		/bin/login				busybox												777		0		0
slink		/bin/ls					busybox												777		0		0
slink		/bin/mkdir				busybox												777		0		0
slink		/bin/mount				busybox												777		0		0
slink		/bin/sh					busybox												777		0		0
slink		/bin/id					busybox												777		0		0

#nod		<name>																			<mode>	<uid>	<gid>	<dev_type>	<maj>	<min>
nod		/dev/console																	0644		0		0		c				5		1

#pipe		<name>																			<mode>	<uid>	<gid>

#sock		<name>																			<mode>	<uid>	<gid>

