# rootfs.list
#dir		<name>																			<mode>	<uid>	<gid>
dir		/sbin																				0755		0		0
dir		/bin																				0755		0		0
	
#file		<name>					<location>											<mode>	<uid>	<gid>	[<hard links>]
file		/bin/busybox			/tks/dctab/br_out-skel/bin/busybox			4755		0		0
file		/etc/inittab			/tks/dctab/etc/initttab							0644		0		0
file		/etc/fstab				/tks/dctab/etc/fstab								0644		0		0
file		/etc/init.d/S01dctab	/tks/dctab/etc/init.d/S01dctab				0755		0		0

#slink	<name>					<target>												<mode>	<uid>	<gid>
slink		/sbin/fsck				../bin/busybox										777		0		0
slink		/sbin/getty		   	../bin/busybox										777		0		0
slink		/sbin/init				../bin/busybox										777		0		0

slink		/sbin/add-shell		../usr/sbin/add-shell							777		0		0

slink		/sbin/blkid				../usr/sbin/blkid									777		0		0
slink		/sbin/depmod			../usr/bin/kmod									777		0		0
slink		/sbin/hwclock			../usr/sbin/hwclock								777		0		0
slink		/sbin/insmod			../usr/bin/kmod									777		0		0
slink		/sbin/ip					../usr/sbin/ip										777		0		0
slink		/sbin/lsmod				../usr/bin/kmod									777		0		0
slink		/sbin/modinfo			../usr/bin/kmod									777		0		0
slink		/sbin/modprobe			../usr/bin/kmod									777		0		0
slink		/sbin/rmmod				../usr/bin/kmod									777		0		0
slink		/sbin/sysctl			../usr/sbin/sysctl								777		0		0


slink		/bin/ash					busybox												777		0		0
slink		/bin/ln					busybox												777		0		0
slink		/bin/login				busybox												777		0		0
slink		/bin/ls					busybox												777		0		0
slink		/bin/mkdir				busybox												777		0		0
slink		/bin/mount				busybox												777		0		0
slink		/bin/sh					busybox												777		0		0
slink		/bin/id					busybox												777		0		0

slink		/bin/bash				../usr/bin/bash									777		0		0
slink		/bin/cat					../usr/bin/coreutils								777		0		0
slink		/bin/chmod				../usr/bin/coreutils								777		0		0
slink		/bin/chown				../usr/bin/coreutils								777		0		0
slink		/bin/cp					../usr/bin/coreutils								777		0		0
slink		/bin/date				../usr/bin/coreutils								777		0		0
slink		/bin/dd					../usr/bin/coreutils								777		0		0
slink		/bin/df					../usr/bin/coreutils								777		0		0
slink		/bin/dmesg				../usr/bin/dmesg									777		0		0
slink		/bin/echo				../usr/bin/coreutils								777		0		0
slink		/bin/false				../usr/bin/coreutils								777		0		0
slink		/bin/free				../usr/bin/free									777		0		0
slink		/bin/kill				../usr/bin/coreutils								777		0		0
slink		/bin/lsblk				../usr/bin/lsblk									777		0		0
slink		/bin/mknod				../usr/bin/coreutils								777		0		0
slink		/bin/mktemp				../usr/bin/coreutils								777		0		0
slink		/bin/mv					../usr/bin/coreutils								777		0		0
slink		/bin/printenv			../usr/bin/coreutils								777		0		0
slink		/bin/ps					../usr/bin/coreutils								777		0		0
slink		/bin/pwd					../usr/bin/coreutils								777		0		0
slink		/bin/rm					../usr/bin/coreutils								777		0		0
slink		/bin/sync				../usr/bin/coreutils								777		0		0
slink		/bin/touch				../usr/bin/coreutils								777		0		0
slink		/bin/true				../usr/bin/coreutils								777		0		0
slink		/bin/umount				../usr/bin/umount									777		0		0
slink		/bin/uname				../usr/bin/coreutils								777		0		0
slink		/bin/uptime				../usr/bin/uptime									777		0		0



#nod		<name>																			<mode>	<uid>	<gid>	<dev_type>	<maj>	<min>
nod		/dev/console																	0644		0		0		c				5		1

#pipe		<name>																			<mode>	<uid>	<gid>

#sock		<name>																			<mode>	<uid>	<gid>

